# debian-control/postinst.in
#!/usr/bin/env bash

set +e

useradd -r lumexcurrency
groupadd -r lumexcurrency

mkdir -p /var/lumexcurrency/Lumex
mkdir -p /var/lumexcurrency/LumexBeta
mkdir -p /var/lumexcurrency/LumexTest
chown -R lumexcurrency: /var/lumexcurrency

systemctl stop ${LUMEX_SERVICE} 2>/dev/null
systemctl disable ${LUMEX_SERVICE} 2>/dev/null
cp ${CMAKE_INSTALL_PREFIX}/extras/systemd/${LUMEX_SERVICE} /etc/systemd/system/${LUMEX_SERVICE}
systemctl daemon-reload
systemctl enable ${LUMEX_SERVICE}
systemctl start ${LUMEX_SERVICE}
sleep 1
systemctl --no-pager status ${LUMEX_SERVICE}

set -e
